# 99-power.rules
SUBSYSTEM=="power_supply", ACTION=="change", RUN+="/home/user/Documents/Scripts/cpupower-control.sh"